# workbench/deploy/backup/workbench-backup.service
[Unit]
Description=Workbench SQLite → GCS backup (one-shot)
Documentation=https://github.com/tripplemay/trade/blob/main/docs/specs/B021-cloud-deploy-auth-spec.md
After=network-online.target workbench-backend.service
Wants=network-online.target
# Backup is allowed to run while workbench-backend is up — sqlite3 .backup
# is a WAL-safe online snapshot.

[Service]
Type=oneshot
User=deploy
Group=deploy
EnvironmentFile=-/etc/workbench/workbench.env
ExecStart=/srv/workbench/current/deploy/backup/workbench-backup.sh
TimeoutStartSec=15min

# The script writes only to /tmp, /var/log/workbench, and GCS. No new
# privileges required.
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
