# workbench/deploy/systemd/workbench-market-context.service
[Unit]
Description=Workbench market-context daily fetch (FRED + Alpha Vantage, read-only)
Documentation=https://github.com/tripplemay/trade/blob/main/docs/specs/B035-market-context-spec.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=deploy
Group=deploy
WorkingDirectory=/srv/workbench/current/backend

# Secrets (FRED_API_KEY, ALPHAVANTAGE_API_KEY, WORKBENCH_DB_URL, ...) live in
# /etc/workbench/workbench.env (chmod 600 root:deploy). Never committed.
EnvironmentFile=/etc/workbench/workbench.env

# Persistent snapshot root so raw provider responses survive release swaps +
# the release GC (mirrors the B033 news snapshot location). Not a secret —
# a fixed path, so it lives here rather than in the env file.
Environment=WORKBENCH_MARKET_SNAPSHOT_DIR=/var/lib/workbench/data/snapshots/market-context

# Permanent product boundary (r) (B035 spec §3): this timer does READ-ONLY
# market-data fetch ONLY. ExecStart invokes the market-context CLI and
# nothing else — it never runs broker / order-ticket / execution /
# recommendation / LLM code. tests/safety/test_market_scheduler_scope.py
# greps this ExecStart line to enforce that.
ExecStart=/opt/workbench/.venv/bin/python -m workbench_api.market.cli fetch

# Standard hardening (mirrors workbench-backend.service): nothing here needs
# a TTY or write access outside the persistent data + log roots.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/workbench /var/log/workbench /tmp

# No [Install] section — this oneshot is pulled by
# workbench-market-context.timer, never enabled / started directly.
