# workbench/deploy/systemd/workbench-data-refresh.service
[Unit]
Description=Workbench real-data refresh (Tiingo prices + SEC EDGAR fundamentals, read-only)
Documentation=https://github.com/tripplemay/trade/blob/main/docs/specs/B045-real-data-refresh-pipeline-spec.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=deploy
Group=deploy
WorkingDirectory=/srv/workbench/current/backend

# B078 F001 — hard wall-clock watchdog. Type=oneshot DISABLES TimeoutStartSec by
# default (it defaults to infinity for oneshot), which is exactly why a single
# hung akshare A-share fetch left this service stuck "activating" for 3 days
# (2026-06-22 → -25), freezing every downstream daily refresh (A-share prices /
# universe / recommendations all stalled on 06-22). Bound the whole run: a normal
# wide refresh is ~33min (B075 VM-measured) plus the per-call fetch timeouts, so
# 90min sits comfortably above the healthy path yet kills a true hang, letting
# systemd fire the next day's run. Paired with --cn-fetch-timeout-seconds below
# (the per-symbol bound) — belt and braces.
TimeoutStartSec=5400

# Secrets (TIINGO_API_KEY, SEC_EDGAR_CONTACT_EMAIL, ...) live in
# /etc/workbench/workbench.env (chmod 600 root:deploy). Never committed.
EnvironmentFile=/etc/workbench/workbench.env

# The VM data store the unified CSVs land under (and that B045 F002 points the
# trade loaders at via WORKBENCH_DATA_ROOT). Persists across release swaps,
# next to the SQLite DB. Not a secret — a fixed path.
Environment=WORKBENCH_DATA_ROOT=/var/lib/workbench/data

# Permanent product boundary (r) (B045 spec §3): this timer does READ-ONLY
# market-data fetch ONLY — Tiingo daily prices + SEC EDGAR quarterly
# fundamentals for the Master universe. ExecStart invokes the data-refresh CLI
# and nothing else — it never runs broker / order-ticket / execution / fills /
# reconcile code. tests/safety/test_market_scheduler_scope.py greps this
# ExecStart line and scans the data_refresh package to enforce that.
#
# B075 F002 — the DAILY job carries the wide A-share universe PRICES (~1500
# liquid names, the §23 VM-reachable sina bulk discovery; the daily 命门 the
# cn_attack momentum scoring reads), but DECOUPLES the two heavy, low-churn cost
# centres: the historical-mcap PIT universe build (~106min @ N=1500) and the CAS
# fundamentals (~116min) run on the separate weekly workbench-cn-universe timer
# (--no-cn-universe-build --no-cn-fundamentals here). VM feasibility GO @ N=1500
# (docs/test-reports/B075-F001-vm-feasibility-probe-2026-06-22.json): daily wide
# prices ~33min, 100% success.
ExecStart=/opt/workbench/.venv/bin/python -m workbench_api.data_refresh.cli fetch --cn-universe-sina-fallback --cn-universe-top-n 1500 --cn-universe-max-superset 1500 --no-cn-universe-build --no-cn-fundamentals --cn-fetch-timeout-seconds 60

# Standard hardening (mirrors workbench-market-context.service).
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/workbench /var/log/workbench /tmp

# No [Install] section — this oneshot is pulled by workbench-data-refresh.timer,
# never enabled / started directly.
