# workbench/deploy/systemd/workbench-recommendations.service
[Unit]
Description=Workbench recommendations precompute (Master Portfolio real scoring, read-only)
Documentation=https://github.com/tripplemay/trade/blob/main/docs/specs/B044-real-scoring-precompute-spec.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=deploy
Group=deploy
WorkingDirectory=/srv/workbench/current/backend

# Secrets (WORKBENCH_DB_URL, ...) live in /etc/workbench/workbench.env
# (chmod 600 root:deploy). Never committed.
EnvironmentFile=/etc/workbench/workbench.env

# B045 F002 — point the trade loaders at the VM data root the B045
# data-refresh job writes (same path workbench-data-refresh.service sets).
# The trade unified-CSV loaders read <root>/snapshots/.../unified/* when this
# is set, so the Master scoring reads REAL refreshed data instead of the
# wheel-bundled B025 fixture. Not a secret — a fixed path. Local / CI leave it
# unset and keep the repo-root behaviour unchanged.
Environment=WORKBENCH_DATA_ROOT=/var/lib/workbench/data

# Permanent product boundary (r-c) (B044 spec §3): this timer runs the
# DETERMINISTIC QUANT scoring precompute ONLY — it imports the trade package
# (Master Portfolio scoring) and writes the recommendation_snapshot table.
# ExecStart invokes the recommendations precompute CLI and nothing else — it
# never runs broker / order-ticket / execution / fills / reconcile code.
# tests/safety/test_market_scheduler_scope.py greps this ExecStart line and
# scans the precompute package to enforce that. Scoring is read-only: it
# never places an order (positioning §1.1 — target weights, not predictions).
ExecStart=/opt/workbench/.venv/bin/python -m workbench_api.recommendations.cli

# Standard hardening (mirrors workbench-market-context.service).
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/workbench /var/log/workbench /tmp

# No [Install] section — this oneshot is pulled by
# workbench-recommendations.timer, never enabled / started directly.
